# .config/systemd/user/mbsync.service
[Unit]
Description=mbsync for GMail

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync gmail

[Install]
WantedBy=default.target
