# .config/systemd/user/notmuch.service
[Unit]
Description=Sync notmuch database

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new

[Install]
WantedBy=default.target
