# .config/systemd/user/rss2email.service
[Unit]
Description=Sync notmuch database

[Service]
Type=oneshot
ExecStart=/usr/bin/r2e run

[Install]
WantedBy=default.target
