# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Sync notmuch database

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync

[Install]
WantedBy=default.target
